# recipes-application/initscript/files/qtapp.service
[Unit]
Description=Main Qt application

[Service]
Environment=QT_QPA_EGLFS_HIDECURSOR=1
ExecStart=nftminter

[Install]
WantedBy=multi-user.target
